# roles/camera-import/files/blinkt-disk-usage.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/blinkt-disk-usage

[Install]
WantedBy=multi-user.target
