# roles/camera-import/files/blinkt-clear.service
[Unit]
DefaultDependencies=no

[Service]
ExecStart=/usr/local/bin/blinkt-clear
Type=oneshot

[Install]
WantedBy=shutdown.target
